# ansible-automation/roles/prometheus_exporters/templates/blackbox_exporter.service.j2
[Unit]
Description=Blackbox Exporter
After=network.target

[Service]
ExecStart=/opt/exporters/blackbox_exporter
Restart=always

[Install]
WantedBy=multi-user.target
