# ansible-automation/roles/prometheus_exporters/templates/nginx_exporter.service.j2
[Unit]
Description=NGINX Exporter
After=network.target

[Service]
ExecStart=/opt/exporters/nginx-prometheus-exporter
Restart=always

[Install]
WantedBy=multi-user.target
